# list-updated-packages/update-list-updated-packages.timer
[Unit]
Description=List updates are available
After=network.target

[Timer]
OnCalendar=*-*-* *:0/10:00

[Install]
WantedBy=timers.target
